# conf/ddnsdpp.service
[Unit]
Description=DDNSD++ (DDNS Daemon)
After=network.target

[Service]
ExecStart=/usr/bin/ddnsd++

[Install]
WantedBy=multi-user.target
